# configs/cron.service
[Service]
Exec=/sbin/bh-crond
Respawn=true
After=network.service
User=root
Group=root
