# configs/network.service
[Service]
Exec=sh -c "ifconfig lo 127.0.0.1 up && ifconfig eth0 up && dhcp -i eth0 -s /bin/dhcp.script"
Respawn=false
